# tsdfileapi/config/file-api.service
[Unit]
Description=file api REST Service

[Service]
User=p01-fileapi-user
Group=p01-fileapi-user-group
ExecStart=/opt/tsd-file-api-venv/virtualenv/bin/python /lib/python2.7/site-packages/tsdfileapi/api.py /etc/tsd-file-api/file-api-config.yaml

[Install]
WantedBy=multi-user.target
